# src/scripts/arch_update/v2.3/archupd.service
# /etc/systemd/system/archupd.service
[Unit]
Description=Arch Linux Auto Update
After=network-online.target
Wants=network-online.target
ConditionACPower=true
ConditionPathExists=/usr/local/bin/archupd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/archupd --auto
# Лимиты для безопасности
MemoryMax=500M
CPUQuota=50%
# Логирование в journal
StandardOutput=journal
StandardError=journal
# Безопасность
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/cache/pacman/pkg /var/lib/pacman

[Install]
WantedBy=multi-user.target